# deploy/trader.service
[Unit]
Description=Trader FX daemon — REST API + live journal
Documentation=https://github.com/rustyeddy/trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/var/lib/trader

# Token lives here so it stays out of the config file and command line.
EnvironmentFile=-/etc/trader/env

ExecStart=/usr/local/bin/trader serve \
    --config /etc/trader/trader.yaml \
    --log-level info

Restart=on-failure
RestartSec=10
TimeoutStopSec=30

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/trader /var/log/trader

StandardOutput=journal
StandardError=journal
SyslogIdentifier=trader

[Install]
WantedBy=multi-user.target
